# node.ini
[program:node]
command=node /home/ltdstats/nodejs/
